# recipes-misc/de10-nano-fpga-init/files/de10-nano-fpga-init.service
[Unit]
Description=Start FPGA init

[Service]
Type=simple
ExecStart=/usr/bin/de10-nano-fpga-init.sh

[Install]
WantedBy=basic.target
